# scripts/cowrie-event-indexer.service
[Unit]
Description=Cowrie Event Indexer Daemon
Documentation=https://github.com/cowrie/cowrie
After=docker.service
Requires=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/cowrie
Environment="PATH=/root/.local/bin:/usr/local/bin:/usr/bin:/bin"

# VT_API_KEY is loaded from /opt/cowrie/etc/report.env by the Python script
# This enables VirusTotal scanning in the event-indexer

# Run the event indexer with container paths
ExecStart=UV_PATH_PLACEHOLDER run /opt/cowrie/scripts/event-indexer.py \
    --db /var/lib/docker/volumes/cowrie-var/_data/lib/cowrie/cowrie.db \
    --log /var/lib/docker/volumes/cowrie-var/_data/log/cowrie/cowrie.json \
    --log-dir /var/lib/docker/volumes/cowrie-var/_data/log/cowrie

# Restart policy
Restart=always
RestartSec=10

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cowrie-event-indexer

[Install]
WantedBy=multi-user.target
